# files/root/service/burp/run
#!/bin/sh
exec /usr/local/sbin/burp -c /etc/burp/burp-server.conf 2>&1
